# etc/systemd/system/kestrel-WebApp.service
[Unit]
Description=Example .NET Web API App running on Ubuntu

[Service]
WorkingDirectory=/var/www/WebApp/bin/Release/netcoreapp3.1/publish/
ExecStart=/usr/bin/dotnet /var/www/WebApp/bin/Release/netcoreapp3.1/publish/WebApp.dll
Restart=always
# Restart service after 10 seconds if the dotnet service crashes:
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=dotnet-example
User=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
